# xendor/alpine-minirootfs-3.24.1-riscv64/etc/inittab
# /etc/inittab for the initramfs.
#
# This image intentionally does not install OpenRC. BusyBox init is PID 1:
# it starts a few one-shot services and respawns serial shells.

::sysinit:/sbin/syslogd
::once:/bin/sh -c '/usr/sbin/sshd -t && /usr/sbin/sshd || echo "warning: sshd failed to start" >/dev/console'

# Respawn login-free serial consoles. BusyBox init opens /dev/<id> for these
# entries, so the shell is attached to the serial device without being PID 1.
ttyS0::respawn:/bin/sh -l
ttyS1::respawn:/bin/sh -l

::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/bin/sync
